# SystemD/.config/systemd/user/sxhkd.service
[Unit]
Description=sxhkd, simple X Hotkey Daemon for use with BSPWM

[Service]
ExecStart=/usr/bin/sxhkd


[Install]
WantedBy=default.target
